# overlay_fs/etc/udev/rules.d/96-autohandleusb.rules
#Handle USB drives by creating a /dev/usbstick symlink and
#executing the bluetracker-usb script.
BUS=="usb", KERNEL=="sd?1", SYMLINK+="usbstick", RUN+="/usr/local/sbin/bluetracker-usb.sh &> /dev/null"
