# overlay_fs/etc/cron.d/voyage
# Default cron file.

# Backup to permanent USB drive every 20 minutes (at *:10 *:30 and *:50)
10,30,50 * * * * root /usr/local/sbin/gyrid-usb.sh cron &> /dev/null

# Call the clean script to free space hourly at *:15
15 * * * * root /usr/local/sbin/gyrid-clean.sh &> /dev/null

# Call the restart script every minute (only restart if necessary)
* * * * * root /usr/local/sbin/gyrid-restart.py &> /dev/null

# Call the NTP restart script every 15 minutes (only restart if necessary)
3,18,33,48 * * * * root /usr/local/sbin/ntp-restart &> /dev/null

# Call the mobile-internet restart script every 15 minutes (only restart if necessary)
8,23,38,53 * * * * root /usr/local/sbin/mobile-internet-restart &> /dev/null
